# linux/tsync.service
[Unit]
Description=tsync S3 FUSE filesystem
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=%h/.local/bin/tsync start
ExecStop=%h/.local/bin/tsync stop
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=default.target
